# crontab
0 */2 * * * /usr/local/bin/mlmmj-maintd -F -d /var/spool/mlmmj
